# crontab
0 1 * * SUN /home/downloader/download-ikea-firmware.sh
